# groovescanner__supervisord.conf
[program:groovescanner]
command=uwsgi -i /home/niespodd/groovescanner/backend/app.ini
